# keycloak.service.j2
[Unit]
Description=Keycloak
After=network.target

[Service]
ExecStart=/opt/keycloak/bin/kc.sh -b=0.0.0.0
User=keycloak
Type=simple

[Install]
WantedBy=multi-user.target
